# tools/git_sync_mds/git_sync_mds.service
[Unit]
Description=MDS Git Synchronization Service
After=wifi-manager.service network-online.target
Wants=network-online.target
# Restart limits (must be in [Unit] for systemd 249+)
StartLimitIntervalSec=600
StartLimitBurst=3

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/home/droneshow/mavsdk_drone_show/tools/update_repo_ssh.sh
WorkingDirectory=/home/droneshow/mavsdk_drone_show
User=droneshow
Group=droneshow
TimeoutStartSec=600
Environment="HOME=/home/droneshow"
Environment="PATH=/home/droneshow/mavsdk_drone_show/venv/bin:/usr/local/bin:/usr/bin:/bin"

# Security hardening
PrivateTmp=yes
NoNewPrivileges=yes

StandardOutput=journal
StandardError=journal
SyslogIdentifier=git-sync-mds

[Install]
WantedBy=multi-user.target
