# tools/wifi-manager/wifi-manager.service
[Unit]
Description=MDS WiFi Manager for Drone Show
# WiFi manager ESTABLISHES network, so it runs before network-online
After=basic.target
Before=network-online.target
# Removed: Wants=network-online.target (contradictory - we run BEFORE it)
# Restart limits (must be in [Unit] for systemd 249+)
StartLimitIntervalSec=300
StartLimitBurst=5

[Service]
Type=simple
ExecStart=/home/droneshow/mavsdk_drone_show/tools/wifi-manager/wifi-manager.sh
WorkingDirectory=/home/droneshow/mavsdk_drone_show/tools/wifi-manager

# Keep as root for WiFi management (nmcli requires elevated privileges)
User=root

# Environment variables for root user
Environment="PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"
Environment="HOME=/root"

# Timeout to prevent boot hang if WiFi takes too long (increased for slow networks)
TimeoutStartSec=180

Restart=on-failure
RestartSec=10

# Security hardening (limited - root needed for nmcli)
PrivateTmp=yes

StandardOutput=journal
StandardError=journal
SyslogIdentifier=wifi-manager

[Install]
WantedBy=multi-user.target
